# home/.xmonad/99-Fujitsu-dock.rules
KERNEL=="dock.1", ACTION=="change", RUN+="/root/Fujitsu_dock.sh >> /root/Fujitsu_dock.log"
